# collector/systemd.service
[Unit]
Description=Caph data collector
After=caph_db.service

[Service]
Type=simple
Environment=MORGAN_COLOR=false
Environment=RUST_BACKTRACE=1
ExecStartPre=/bin/sleep 2
ExecStart=/opt/caph/caph_collector
Restart=always

[Install]
WantedBy=multi-user.target
